# srcs/grafana/supervisor.conf
[supervisord]
nodaemon = true

[program: grafana]
command = sh -c "/usr/local/bin/grafana-server -homepath /grafana && kill -s SIGTERM $(cat supervisord.pid)"