# srcs/ftps/supervisor.conf
[supervisord]
nodaemon = true

[program:vsftpd]
command = sh -c "/usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf && kill -s SIGTERM $(cat supervisord.pid)"